# ffffm-button-bind/files/etc/init.d/button-bind
#!/bin/sh /etc/rc.common

FUNCTION=$(uci get button-bind.wifi.function)
if [ "${FUNCTION}" = 3 ]
then
	for i in /sys/class/leds/* ; do echo 0 > "$i"/brightness ; done
fi

exit 0
